# templates/default/metrics.cron.erb
*/5 * * * * root <%= @chatmail_bin %>/chatmail-metrics <%= @config['mailboxes_dir'] %>/<%= @config['domain'] %> >/var/www/html/metrics
