# templates/debian/filtermail.service.erb
[Unit]
Description=Outgoing Chatmail Postfix before queue filter

[Service]
<% unless @dkim_enable %>
Environment="FILTERMAIL_SKIP_DKIM=1"
<% end %>
ExecStart=<%= @execpath %> <%= @config_path %> outgoing
Restart=always
RestartSec=30
User=vmail

[Install]
WantedBy=multi-user.target
