# templates/debian/lego-renewal.service.erb
[Unit]
Description=LetsEncrypt certificate renewal

[Service]
Type=oneshot
<% @lego_dns_envs.each do |key, value| %>
Environment="<%= key %>=<%= value %>"
<% end %>
ExecStart=<%= @lego_issuance_command %> --deploy-hook=<%= @lego_deploy_hook %>

[Install]
WantedBy=multi-user.target
