# templates/freebsd/doveauth.erb
#!/bin/sh
#
# PROVIDE: doveauth
# REQUIRE: LOGIN
# BEFORE: dovecot
# KEYWORD: shutdown

. /etc/rc.subr

name="doveauth"
rcvar="doveauth_enable"

load_rc_config $name

: ${doveauth_enable:="NO"}
: ${doveauth_svc_user:="vmail"}
: ${chatmail_config:="<%= @config_path %>"}

required_files="${chatmail_config}"
procname="*python*"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args=" -T ${name} -p ${pidfile} -u ${doveauth_svc_user} -f <%= @execpath %> /var/run/${name}/${name}.sock ${chatmail_config}"

start_precmd="install -d -o ${doveauth_svc_user} -g ${doveauth_svc_user} -m 755 /var/run/${name}"

run_rc_command "$1"
