# templates/freebsd/filtermail.erb
#!/bin/sh
#
# PROVIDE: filtermail
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="filtermail"
rcvar="filtermail_enable"

load_rc_config $name

: ${filtermail_enable:="NO"}
: ${filtermail_user:="vmail"}
: ${filtermail_execpath:="<%= @execpath %>"}
: ${filtermail_config_path:="<%= @config_path %>"}

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args=" -T ${name} -P ${pidfile} -u ${filtermail_user} -r -f ${filtermail_execpath} ${filtermail_config_path} outgoing"

start_precmd="filtermail_prestart"
start_cmd="filtermail_start"
stop_cmd="filtermail_stop"

filtermail_prestart()
{
    if [ ! -f "${filtermail_config_path}" ]; then
        echo "Config file ${filtermail_config_path} does not exist!"
        return 1
    fi
}

filtermail_start()
{
    echo "Starting filtermail service."
    eval ${command} ${command_args}
}

filtermail_stop()
{
    echo "Stopping filtermail service."
    pkill -f "${filtermail_execpath}"
}

run_rc_command "$1"
