# templates/freebsd/filtermail.erb
#!/bin/sh
#
# PROVIDE: filtermail
# REQUIRE: LOGIN
# BEFORE: postfix
# KEYWORD: shutdown

. /etc/rc.subr

name="filtermail"
rcvar="filtermail_enable"

load_rc_config $name

: ${filtermail_enable:="NO"}
: ${filtermail_svc_user:="vmail"}
: ${chatmail_config:="<%= @config_path %>"}

required_files="${chatmail_config}"
procname="*python*"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args=" -T ${name} -P ${pidfile} -u ${filtermail_svc_user} -r -f <%= @execpath %> ${chatmail_config} outgoing"

run_rc_command "$1"
