# templates/freebsd/filtermail-incoming.erb
#!/bin/sh
#
# PROVIDE: filtermail_incoming
# REQUIRE: LOGIN
# BEFORE: postfix
# KEYWORD: shutdown

. /etc/rc.subr

name="filtermail_incoming"
rcvar="filtermail_incoming_enable"

load_rc_config $name

: ${filtermail_incoming_enable:="NO"}
: ${filtermail_svc_user:="vmail"}
: ${chatmail_config:="<%= @config_path %>"}

required_files="${chatmail_config}"
procname="*python*"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args=" -T ${name} -P ${pidfile} -u ${filtermail_svc_user} -r -f <%= @execpath %> ${chatmail_config} incoming"

run_rc_command "$1"
