# templates/freebsd/filtermail.erb
#!/bin/sh
#
# PROVIDE: filtermail
# REQUIRE: LOGIN
# BEFORE: postfix
# KEYWORD: shutdown

. /etc/rc.subr

name="filtermail"
rcvar="filtermail_enable"

load_rc_config $name

: ${filtermail_enable:="NO"}
: ${filtermail_svc_user:="vmail"}
: ${chatmail_config:="<%= @config_path %>"}

required_files="${chatmail_config}"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args=" -T ${name} -p ${pidfile} -u ${filtermail_svc_user} -f <%= @execpath %> ${chatmail_config} outgoing"

# Override /etc/rc.subr variant because we do not
# care about procname matching and we want to be
# able to handle the Python to Rust filtermail
# transition gracefully.
check_pidfile()
{
	if [ -f $pidfile ]; then
	    read _pid _junk < $pidfile
	    echo $_pid
	else
	    echo ""
	fi
}

run_rc_command "$1"
