# templates/freebsd/lastlogin.erb
#!/bin/sh
#
# PROVIDE: lastlogin
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="lastlogin"
rcvar="lastlogin_enable"

load_rc_config $name

: ${lastlogin_enable:="NO"}
: ${lastlogin_svc_user:="vmail"}
: ${chatmail_config:="<%= @config_path %>"}

required_files="${chatmail_config}"
procname="*python*"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args=" -T ${name} -p ${pidfile} -u ${lastlogin_svc_user} -r -f <%= @execpath %> /var/run/${name}/lastlogin.sock ${chatmail_config}"

start_precmd="install -d -o ${lastlogin_svc_user} -g ${lastlogin_svc_user} -m 755 /var/run/${name}"

run_rc_command "$1"
